# overlay/etc/sv/rsyslog/run
#!/bin/sh

set -eu

exec rsyslogd -n -f /etc/rsyslog.conf
